# minecraft.cron
#!/bin/bash
#This script is designed to be a cron that checks to see if the Minecraft server is up.
#If it is not up, the scrip will start it, if it is up, it will do nothing.

if [[ $(echo "`ps auwx |grep 'minecraft_server' |grep -v grep |awk '{print$11}'`") != "java" ]] ;
then cd /home/minecraf/minecraft ; java -Xmx1024M -Xms1024M -jar /home/minecraf/minecraft/minecraft_server.1.8.7.jar nogui
fi
